# containers/web/extensions/zoco/cron.ini
[program:cron]
command = /usr/sbin/crond -f
user = root
autostart = true
